# deploy/ghost-trader.service
[Unit]
Description=Kalshi Ghost Trader Backend
After=network.target

[Service]
Type=simple
User=fahad
WorkingDirectory=/home/fahad/kalshi-ghost-trader
Environment=APP_ENV=prod
ExecStart=/home/fahad/kalshi-ghost-trader/ghost-trader
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
